# templates/certbot-renew.service
[Unit]
Description=Certbot Renew

[Service]
Type=simple
ExecStart=/usr/bin/certbot renew -n
ExecStopPost=-/bin/systemctl reload nginx.service
